# systemd_units/mirr_upd.service
[Unit]
Description=Upgrades Mirrorlist

[Service]
Type=simple
ExecStart=/home/simone/sys_maint/pacman/mirr_upd.sh
