# systemd_units/mirr_upd.timer
[Unit]
Description=mirr_upd.service

[Timer]
Unit=mirr_upd.service
OnBootSec=10

[Install]
WantedBy=timers.target
