# systemd_units/mylogs.service
[Unit]
Description=Saves logs files

[Service]
Type=simple
ExecStart=/home/simone/sys_maint/logs/logs.sh


