# systemd_units/pac_cache.service
[Unit]
Description=Clears cached packages older than three updates ago

[Service]
Type=simple
ExecStart=/home/simone/sys_maint/pacman/rm_pac_cac.sh
