# systemd_units/pac_cache.timer
[Unit]
Description=Clears cached packages older than three upgrades

[Timer]
Unit=pac_cac.service
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
