# systemd_units/rmc.service
[Unit]
Description=Cleans cahce on boot

[Service]
Type=simple
ExecStart=python /home/simone/sys_maint/cleanup/rm.py

